# dist/common-poc/templates/etc/profile.d/harvard.sh.erb
#### THIS FILE MANAGED BY PUPPET ####
alias view="/usr/bin/vim -R"
alias facter="facter -p"

bind '"\e[A"':history-search-backward
bind '"\e[B"':history-search-forward

export http_proxy=http://<%= @proxyserver -%>:3128
export https_proxy=http://<%= @proxyserver -%>:3128
export no_proxy="localhost,127.0.0.1,<%= @ipaddress -%>,.localdomain.com,.<%= @domain -%>"

#### THIS FILE MANAGED BY PUPPET ####

